# DEBIAN/postinst
#!/bin/bash

systemctl enable feren-flatpak-management.service